# laravel-backend/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Laravel scheduler — container-native ticker (replaces a system crontab). Runs
# `schedule:run` every minute internally, dispatching the payment reconciliation
# (every 5 min) and courier status poll (hourly). Ships with the image, so it is
# always on and restarts on any deploy — no separate EasyPanel service needed.
[program:scheduler]
command=php /var/www/html/artisan schedule:work
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Queue worker — processes the queued jobs (payment reconciliation sweep, courier
# auto-push on confirm, status sync). QUEUE_CONNECTION=database. `--max-time`
# recycles the worker hourly to avoid memory creep; `stopwaitsecs` lets an
# in-flight job finish gracefully on redeploy instead of being killed mid-write.
[program:queue-worker]
command=php /var/www/html/artisan queue:work --tries=3 --backoff=60 --sleep=3 --max-time=3600
user=www-data
autostart=true
autorestart=true
priority=20
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
